# config/includes.chroot_after_packages/etc/profile.d/init.sh
#!/bin/bash
path=$(pwd)
cd /hardwarecheck; sudo ./init.sh live; ./inventory.sh; cd $path
